# .config/systemd/user/pipewire-filter-chain.service
[Unit]
Description=Pipewire filter chain

[Service]
ExecStart=pipewire -c filter-chain.conf

[Install]
WantedBy=default.target
